# src/mycat1.6/supervisord.conf
[supervisord]
nodaemon=true
[program:mycat]
command=/usr/local/mycat/bin/mycat start