# src/zookeeper3.4.11/supervisord.conf
[supervisord]
nodaemon=true
[program:zookeeper]
command=/usr/local/zookeeper-3.4.11/bin/zkServer.sh start