# src/kafka/supervisord.conf
[supervisord]
nodaemon=true
[program:kafka]
command=/usr/local/kafka_2.11-1.0.0/bin/kafka-server-start.sh /usr/local/kafka_2.11-1.0.0/config/server.properties