# src/kafka/supervisord.conf
[supervisord]
nodaemon=true
[program:mycat]
command=/usr/local/kafka_2.11-1.0.0/bin/kafka-console-producer.sh /usr/local/kafka_2.11-1.0.0/config/server.properties