# src/activeMQ5.15.2/supervisord.conf
[supervisord]
nodaemon=true
[program:activeMQ]
command=service activemq start