# src/activeMQ5.15.2/supervisord.conf
[supervisord]
nodaemon=true
[program:redis]
command=service activemq start