# src/mysql5.6/supervisord.conf
[supervisord]
nodaemon=true
[program:mysqld]
command=service mysqld start