# src/nexus2.11.2-03/supervisord.conf
[supervisord]
nodaemon=true
[program:nexus]
command=/usr/local/nexus/nexus-2.11.2-03/bin/nexus start
