# src/kong0.11.2.el7.noarch/supervisord.conf
[supervisord]
nodaemon=true
[program:kongMigrationsUp]
command=kong migrations up [-c /etc/kong/kong.conf]
[program:kongStart]
command=kong start [-c /etc/kong/kong.conf]