# src/kong0.11.2.el7.noarch/supervisord.conf
[supervisord]
nodaemon=true
[program:tomcat]
command=kong start [-c /etc/kong/kong.conf]